# static/etc/services.d/headphones/run
#!/usr/bin/with-contenv bash

umask 0002

exec s6-setuidgid bytesized python /app/headphones/Headphones.py --datadir /config
